# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:fastapi]
command=uvicorn nerpblog:app --port 9001 --host 0.0.0.0

[program:aiogram]
command=python nerpblog/bot.py

